# backend/deployment/systemd/system/sshd.service.d/priority.conf
[Service]
Nice=-10 # This configuration sets the priority of the sshd service to a higher level.
